# deploy/cage.service
[Unit]
Description=Family Dashboard Kiosk (cage + Chromium)
After=dashboard.service seatd.service network-online.target
Requires=dashboard.service seatd.service

[Service]
Type=simple
User=dashboard
Group=dashboard
PAMName=login
TTYPath=/dev/tty7
StandardInput=tty
StandardOutput=journal
StandardError=journal
UtmpIdentifier=tty7
UtmpMode=user

# Let systemd manage the runtime dir (auto-created at /run/cage-dashboard with
# 0700 owned by User=). XDG_RUNTIME_DIR is what cage/Wayland use for the
# wayland-0 socket. Avoids the systemd-logind /run/user/$uid path which
# doesn't exist for system users.
RuntimeDirectory=cage-dashboard
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/cage-dashboard
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=MOZ_ENABLE_WAYLAND=1

# libseat: skip logind backend (we're a system service with no session),
# go straight to seatd.
Environment=LIBSEAT_BACKEND=seatd

# wlroots on Pi 4: the V3D DRM driver rejects the buffer formats the GL
# renderer asks for, producing "Swapchain ... failed test" errors. Pixman
# is software-only but more than fast enough for this UI, and the legacy
# KMS path avoids atomic modesetting bugs.
Environment=WLR_RENDERER=pixman
Environment=WLR_DRM_NO_ATOMIC=1
Environment=WLR_NO_HARDWARE_CURSORS=1

# Hide the cursor via a blank xcursor theme. install.sh writes
# /usr/share/icons/blank for us.
Environment=XCURSOR_THEME=blank
Environment=XCURSOR_SIZE=1

# `cage-rotated` is a tiny wrapper that calls wlr-randr after the Wayland
# socket appears, applying $CAGE_TRANSFORM to $CAGE_OUTPUT. Override either
# via a systemd drop-in (e.g. set CAGE_TRANSFORM=normal to disable rotation,
# or =270 to flip the other way).
Environment=CAGE_OUTPUT=HDMI-A-1
Environment=CAGE_TRANSFORM=90
ExecStart=/usr/local/bin/cage-rotated -s -- /usr/bin/chromium \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --check-for-update-interval=31536000 \
  --no-first-run \
  --ozone-platform=wayland \
  --enable-features=UseOzonePlatform \
  --start-fullscreen \
  http://localhost:3000/

Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
