# deploy/cage.service
[Unit]
Description=Family Dashboard Kiosk (cage + Chromium)
After=dashboard.service network-online.target
Requires=dashboard.service

[Service]
Type=simple
User=dashboard
Group=dashboard
PAMName=login
TTYPath=/dev/tty7
StandardInput=tty
StandardOutput=journal
StandardError=journal
UtmpIdentifier=tty7
UtmpMode=user

# Let systemd manage the runtime dir (auto-created at /run/cage-dashboard with
# 0700 owned by User=). XDG_RUNTIME_DIR is what cage/Wayland use for the
# wayland-0 socket. Avoids the systemd-logind /run/user/$uid path which
# doesn't exist for system users.
RuntimeDirectory=cage-dashboard
RuntimeDirectoryMode=0700
Environment=XDG_RUNTIME_DIR=/run/cage-dashboard
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=MOZ_ENABLE_WAYLAND=1

ExecStart=/usr/bin/cage -s -- /usr/bin/chromium \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --check-for-update-interval=31536000 \
  --no-first-run \
  --ozone-platform=wayland \
  --enable-features=UseOzonePlatform \
  --start-fullscreen \
  http://localhost:3000/

Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
