# deploy/dashboard.service
[Unit]
Description=Family Dashboard (Node server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dashboard
Group=dashboard
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=DATA_DIR=/var/lib/dashboard
Environment=LOCAL_PHOTOS_DIR=/var/lib/dashboard/photos
EnvironmentFile=-/etc/dashboard/env
WorkingDirectory=/opt/dashboard
ExecStart=/usr/bin/pnpm --filter @dashboard/server start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome must stay off — corepack-managed pnpm lives in
# /home/dashboard/.cache/node/corepack and the service needs to read it.
ProtectHome=false
PrivateTmp=true
ReadWritePaths=/var/lib/dashboard /home/dashboard/.cache

[Install]
WantedBy=multi-user.target
